# root/etc/services.d/aria2/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid abc aria2c --conf-path=/config/aria2.conf
